# build/energy.service
[Unit]
Description=FastAPI Application for Energy Audit Demo
After=network.target

[Service]
WorkingDirectory=/home/randyzwguo/pdf
Enviroment=PATH=/home/randyzwguo/pdf/venv/bin
ExecStart=/home/randyzwguo/pdf/venv/bin/python main.py --host=127.0.0.1 --port=8000 --reload
Restart=always

[Install]
WantedBy=multi-user.target